# roles/alastria-regular-node/files/geth.service
[Unit]
Description=Ethereum go client
Requires=constellation

[Service]
Type=simple
Restart=always
ExecStart=/root/lacchain/start-quorum.sh

[Install]
WantedBy=default.target
